# .config/systemd/user/backup-from-yt.service
# Author: Rafael Cavalcanti - rafaelc.org

[Unit]
Description=Backup from yt
After=network.target
OnFailure=status-email@%n.service

[Service]
ExecStart=%h/Documents/Sysadmin/Backups/bin/backup-from-yt
Nice=19
